# deploy/ansible/templates/brai-api.service.j2
[Unit]
Description=Brai API {{ item.key }}
After=network.target

[Service]
Type=simple
WorkingDirectory={{ brai_env_root }}/{{ item.value.path }}/source/services/brai_api
{% if item.key == 'prod' %}
EnvironmentFile={{ brai_api_env_source }}
ExecStart=/usr/bin/env PORT={{ item.value.api_port }} BRAI_DB={{ brai_repo }}/data/brai.sqlite BRAI_RELEASE_DIR={{ brai_repo }}/deploy/releases BRAI_INBOUND_STORAGE_ROOT={{ brai_repo }}/data/inbox-attachments BRAI_CODEX_BIN=/srv/opt/codex-cli/bin/codex BRAI_CODEX_MODEL=gpt-5.4-mini BRAI_CODEX_TIMEOUT_MS=20000 CODEX_HOME=/srv/opt/codex-runtime/brai HOME=/srv/opt/codex-runtime/brai {{ brai_node_bin }} src/index.js
{% else %}
EnvironmentFile=-{{ brai_env_root }}/{{ item.value.path }}/brai-api.env
ExecStart=/usr/bin/env PORT={{ item.value.api_port }} BRAI_DB={{ brai_env_root }}/{{ item.value.path }}/data/brai.sqlite BRAI_RELEASE_DIR={{ brai_repo }}/deploy/releases BRAI_INBOUND_STORAGE_ROOT={{ brai_env_root }}/{{ item.value.path }}/data/inbox-attachments BRAI_CODEX_BIN=/srv/opt/codex-cli/bin/codex BRAI_CODEX_MODEL=gpt-5.4-mini BRAI_CODEX_TIMEOUT_MS=20000 CODEX_HOME=/srv/opt/codex-runtime/brai HOME=/srv/opt/codex-runtime/brai {{ brai_node_bin }} src/index.js
{% endif %}
Restart=always
RestartSec=3
User={{ brai_service_user }}
{% if item.key == 'prod' %}
Group={{ brai_service_group }}
{% else %}
Group={{ brai_deploy_user }}
{% endif %}
SupplementaryGroups={{ brai_deploy_user }}
UMask=0002
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
