# deploy/systemd/brai-api.service
[Unit]
Description=Brai API
After=network.target

[Service]
Type=simple
User=brai
Group=brai
SupplementaryGroups=brai-deploy
UMask=0002
WorkingDirectory=/srv/projects/brai/services/brai_api
EnvironmentFile=/etc/brai/brai-api.env
Environment=BRAI_INBOUND_STORAGE_ROOT=/srv/projects/brai/data/inbox-attachments
Environment=BRAI_CODEX_BIN=/srv/opt/codex-cli/bin/codex
Environment=BRAI_CODEX_MODEL=gpt-5.4-mini
Environment=BRAI_CODEX_TIMEOUT_MS=20000
Environment=CODEX_HOME=/srv/opt/codex-runtime/brai
Environment=HOME=/srv/opt/codex-runtime/brai
ExecStart=/srv/opt/node-v22.16.0/bin/node src/index.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
